# php56-apache/supervisord.conf
[supervisord]
nodaemon=true
[program:cron]
command=/etc/init.d/cron start
[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
